# deploy/systemd/inference-miner.service
[Unit]
Description=Reliquary Miner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=reliquary
WorkingDirectory=/srv/reliquary-inference/current
EnvironmentFile=/srv/reliquary-inference/state/reliquary-inference.env
EnvironmentFile=-/srv/reliquary-inference/state/inference-miner.env
ExecStart=/bin/bash -lc 'exec /opt/miniforge/bin/conda run --no-capture-output -n reliquary-inference reliquary-inference run-miner'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
